# geth-docker/resources/supervisord.conf
[supervisord]
nodaemon=true

[program:geth]
priority=1
command=/geth/start.sh
stdout_logfile=/geth/%(program_name)s.log
stderr_logfile=/geth/%(program_name)s.log
autorestart=true
startsecs=0
